# scheduler/autolm.timer
[Unit]
Description=Start Automatic Leaderboard Maker monthly

[Timer]
OnCalendar= *-*-1 00:00:00

[Install]
WantedBy=timers.target